# scripts host/loop_ina219.service
# Locate: /etc/systemd/system


[Unit]
Description=INA219 loop
After=docker.service hassio-supervisor.service

[Service]
User=pi
Group=pi
Type=simple
ExecStart=/etc/systemd/system/loop_ina219.sh

[Install]
WantedBy=multi-user.target